# DEBIAN/prerm
#!/bin/bash
systemctl stop ssp.service
systemctl disable ssp.service